# containers/systemd-templets/sillytavern.container
[Service]
Restart=always
RestartSec=60


[Install]
WantedBy=default.target


[Container]
ContainerName=sillytavern
Image=ghcr.io/sillytavern/sillytavern:latest

PublishPort=6700:8000

Volume=%h/containers/sillytavern/config:/home/node/app/config
Volume=%h/containers/sillytavern/data:/home/node/app/data
Volume=%h/containers/sillytavern/plugins:/home/node/app/plugins
Volume=%h/containers/sillytavern/extensions:/home/node/app/public/scripts/extensions/third-party

Environment=NODE_ENV=production
Environment=FORCE_COLOR=1
Volume=/etc/localtime:/etc/localtime:ro
